# config/supervisord.conf

[supervisord]
nodaemon=true

[program:uwsgi-app]
command = /usr/local/bin/uwsgi --ini /home/config/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx